# huggingface/supervisord.conf
; Runs the whole Sentinel stack inside one Hugging Face Space container.
; Children inherit the environment exported by entrypoint.sh.
[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0

; ── Postgres (app audit log + blocklist) ─────────────────────────────────────
[program:postgres]
command=%(ENV_PGBIN)s/postgres -D %(ENV_PGDATA)s -c listen_addresses=127.0.0.1 -p 5432
user=postgres
priority=10
autorestart=true
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── Redis (rate limiting) ────────────────────────────────────────────────────
[program:redis]
command=redis-server --bind 127.0.0.1 --port 6379 --save ""
priority=10
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── OPA (policy engine) ──────────────────────────────────────────────────────
[program:opa]
command=/usr/local/bin/opa run --server --addr=127.0.0.1:8181 --log-level=error /app/policies
priority=10
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── Keycloak (OAuth 2.1 / OIDC) ──────────────────────────────────────────────
[program:keycloak]
command=/opt/keycloak/bin/kc.sh start-dev --import-realm --http-host=127.0.0.1 --http-port=8080
priority=20
autorestart=true
startsecs=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── Sentinel MCP server ──────────────────────────────────────────────────────
; Retries generously: it waits on Postgres coming up on first boot.
[program:sentinel]
command=python -m sentinel.main
directory=/app
priority=30
autorestart=true
startretries=20
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── Caddy (single public port 7860 → app + keycloak + landing) ───────────────
[program:caddy]
command=/usr/local/bin/caddy run --config /etc/caddy/Caddyfile --adapter caddyfile
priority=40
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
